# systemd/robot1/gp-lidar-idle.service
[Unit]
Description=GP LiDAR idle hold (motor off while the robot stack is down)
ConditionPathExists=/dev/rplidar
# ONE-DIRECTIONAL hand-off (a Conflicts= here is bidirectional and let a
# late-starting holder KILL a freshly restarted stack — observed live):
# gp-robot1's ExecStartPre stops this holder before claiming the port;
# its ExecStopPost starts it again. If this holder ever starts while the
# stack runs, the script sees the port busy and exits 0 harmlessly.

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/GP/tools/lidar_idle.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
