# systemd/robot1/gp-robot1.service
[Unit]
Description=GP Robot1 ROS stack (lidar + slam + bridge + explorer + gateway)
Requires=gp-preflight.service
After=gp-preflight.service network-online.target

[Service]
Type=simple
User=pi
Environment=GP_DIR=/home/pi/GP
WorkingDirectory=/home/pi/GP
# Purge stale FastDDS shared-memory segments before every start — after a
# crash they silently eat all messages (see gp-robot2.service).
ExecStartPre=/bin/sh -c 'rm -f /dev/shm/fastrtps_* /dev/shm/fast_datasharing* /dev/shm/sem.fastrtps_* 2>/dev/null; true'
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && if [ -f "$HOME/ros2_ws/install/setup.bash" ]; then source "$HOME/ros2_ws/install/setup.bash"; fi; exec ros2 launch /home/pi/GP/robots/robot1/launch/robot1.launch.py'
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=10
MemoryMax=2500M

[Install]
WantedBy=multi-user.target
