# systemd/robot1/gp-robot1.service
[Unit]
Description=GP Robot1 ROS stack (lidar + slam + bridge + explorer + gateway)
Requires=gp-preflight.service
After=gp-preflight.service network-online.target

[Service]
Type=simple
User=pi
Environment=GP_DIR=/home/pi/GP
# Deterministic DDS transport: UDP loopback only, SHM disabled — the
# shared-memory transport was a per-boot lottery (see the XML's header)
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/pi/GP/config/fastdds_udp_only.xml
WorkingDirectory=/home/pi/GP
# Purge stale FastDDS shared-memory segments before every start — after a
# crash they silently eat all messages (see gp-robot2.service).
ExecStartPre=/bin/sh -c 'rm -f /dev/shm/fastrtps_* /dev/shm/fast_datasharing* /dev/shm/sem.fastrtps_* 2>/dev/null; true'
# Evict the lidar idle holder so the driver can claim the serial port
# ('+' = run as root; one-directional — see gp-lidar-idle.service)
ExecStartPre=+/bin/systemctl stop gp-lidar-idle.service
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && if [ -f "$HOME/ros2_ws/install/setup.bash" ]; then source "$HOME/ros2_ws/install/setup.bash"; fi; exec ros2 launch /home/pi/GP/robots/robot1/launch/robot1.launch.py'
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=10
MemoryMax=2500M
# When the stack stops, hand the lidar back to the idle holder so the
# motor stops spinning ('+' = run as root regardless of User=)
ExecStopPost=+/bin/systemctl --no-block start gp-lidar-idle.service

[Install]
WantedBy=multi-user.target
