# systemd/robot2/gp-robot2.service
[Unit]
Description=GP Robot2 ROS stack (bridge + odom + goto + gateway)
Requires=gp-preflight.service
After=gp-preflight.service network-online.target

[Service]
Type=simple
User=pi
Environment=GP_DIR=/home/pi/GP
# Deterministic DDS transport: UDP loopback only, SHM disabled — the
# shared-memory transport was a per-boot lottery (see the XML's header)
Environment=FASTRTPS_DEFAULT_PROFILES_FILE=/home/pi/GP/config/fastdds_udp_only.xml
WorkingDirectory=/home/pi/GP
# Brownout crash-loops leave stale FastDDS shared-memory segments in
# /dev/shm; with ROS_LOCALHOST_ONLY=1 they silently eat every message
# (discovery matches, data never delivered — field case 2026-06-11:
# bridge executor idle while 695 gateway commands vanished). Purge them
# before every start; nothing else may own them at this point.
ExecStartPre=/bin/sh -c 'rm -f /dev/shm/fastrtps_* /dev/shm/fast_datasharing* /dev/shm/sem.fastrtps_* 2>/dev/null; true'
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && exec ros2 launch /home/pi/GP/robots/robot2/launch/robot2.launch.py'
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=10
MemoryMax=700M

[Install]
WantedBy=multi-user.target
