# systemd/robot2/gp-robot2.service
[Unit]
Description=GP Robot2 ROS stack (bridge + odom + goto + gateway)
Requires=gp-preflight.service
After=gp-preflight.service network-online.target

[Service]
Type=simple
User=pi
Environment=GP_DIR=/home/pi/GP
WorkingDirectory=/home/pi/GP
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && exec ros2 launch /home/pi/GP/robots/robot2/launch/robot2.launch.py'
Restart=always
RestartSec=3
KillMode=mixed
TimeoutStopSec=10
MemoryMax=700M

[Install]
WantedBy=multi-user.target
